# gunicorn.conf
[program:gunicorn]
command=gunicorn --worker-class eventlet -w 1 karel-arena:app --error-logfile /var/log/gunicorn.log --reload
directory=/app
